# LFTR-Broadcast-GFS-Marine-Globe-codex-apply-patch-to-installer-and-clean-directory/deploy/systemd/broadcast.service
[Unit]
Description=LFTR Broadcast / GFS Marine Globe
After=network-online.target
Wants=network-online.target

[Service]
Environment=APP_DIR=${APP_DIR}
Environment=VENV_DIR=${APP_DIR}/venv
Environment=APP_USER=${APP_USER}
Environment=APP_GROUP=${APP_GROUP}
Type=simple
User=${APP_USER}
Group=${APP_GROUP}
WorkingDirectory=${APP_DIR}
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=MALLOC_ARENA_MAX=2
Environment=OMP_NUM_THREADS=1
Environment=OPENBLAS_NUM_THREADS=1
Environment=NUMEXPR_NUM_THREADS=1
Environment=PYTHONMALLOC=malloc
Environment=HDF5_USE_FILE_LOCKING=FALSE
Environment=GFS_TILE_WARM_WORKERS=1
Environment=GFS_TILE_READ_WORKERS=1
Environment=GFS_ALLOW_SYNTHETIC_FALLBACK=0
Environment=ALLOW_SYNTHETIC_FALLBACK=0
EnvironmentFile=/etc/broadcast/install.env
ExecStartPre=/bin/bash -lc 'test -x ${APP_DIR}/venv/bin/python && ${APP_DIR}/venv/bin/python -c "import hypercorn, quart"'
ExecStart=${APP_DIR}/scripts/run_broadcast_service.sh
Restart=always
RestartSec=3
KillMode=mixed
TimeoutStopSec=25
LimitNOFILE=20000
MemoryMax=10G

[Install]
WantedBy=multi-user.target
