# LFTR-Broadcast-GFS-Marine-Globe-codex-apply-patch-to-installer-and-clean-directory/deploy/templates/app.service.template
[Unit]
Description=LFTR Broadcast / GFS Marine Globe
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__
Environment=APP_DIR=__APP_DIR__
Environment=VENV_DIR=__APP_DIR__/venv
Environment=APP_USER=__APP_USER__
Environment=APP_GROUP=__APP_GROUP__
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=MALLOC_ARENA_MAX=2
Environment=OMP_NUM_THREADS=1
Environment=OPENBLAS_NUM_THREADS=1
Environment=NUMEXPR_NUM_THREADS=1
Environment=PYTHONMALLOC=malloc
Environment=HDF5_USE_FILE_LOCKING=FALSE
Environment=GFS_TILE_WARM_WORKERS=1
Environment=GFS_TILE_READ_WORKERS=1
Environment=GFS_ALLOW_SYNTHETIC_FALLBACK=0
Environment=ALLOW_SYNTHETIC_FALLBACK=0
EnvironmentFile=/etc/broadcast/install.env
ExecStartPre=/bin/bash -lc 'test -x __APP_DIR__/venv/bin/python && __APP_DIR__/venv/bin/python -c "import hypercorn, quart"'
ExecStart=__APP_DIR__/scripts/run_broadcast_service.sh
Restart=always
RestartSec=3
KillMode=mixed
TimeoutStopSec=25
LimitNOFILE=20000
MemoryMax=10G

[Install]
WantedBy=multi-user.target
